# yggdrasil-stop.path.in
# This file is part of yggdrasil.
#
# Any changes made to this file will be overwritten during a software update. To
# override a parameter in this file, create a drop-in file, typically located at
# /etc/systemd/system/yggdrasil-stop.path.d/override.conf. Put the desired
# overrides in that file and reload systemd.
#
# For more information about systemd drop-in files, see systemd.unit(5).

[Unit]
Description=Automatically stop yggdrasil by watching if rhsm cert exists
Documentation=man:yggd(1)

[Path]
PathModified=@sysconfdir@/pki/consumer/cert.pem

[Install]
WantedBy=multi-user.target
